# systemd/system/yggdrasil.service.d/override.conf.tmpl
{{ if lookPath "yggdrasilctl" -}}
[Service]
StandardOutput=append:/tmp/yggdrasil.log
{{- end -}}
