# systemd/system/ly.service.d/override.conf.tmpl
{{ if lookPath "ly-dm" -}}
[Service]
Restart=always
RestartSec=1s
{{- end -}}
