# systemd/system/atuin@.service.d/override.conf.tmpl
{{ if eq (output "sh" "-c" "systemctl is-enabled atuin@proteus.service || true" | trim) "enabled" -}}
[Unit]
BindsTo=mnt-storage3.mount
After=mnt-storage3.mount
{{ end -}}
