# systemd/system/lighttpd.socket.d/override.conf.tmpl
{{ if lookPath "lighttpd" -}}
[Socket]
ListenStream=443
{{ end -}}
